# luci-theme-neobird/root/etc/uci-defaults/30_luci-theme-neobird
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Neobird=/luci-static/neobird-pc
	set luci.main.mediaurlbase=/luci-static/neobird-pc
	commit luci
EOF
exit 0 
# ↑↑↑ DO NOT EDIT THIS LINE
