# luci-app-oscam/files/oscam.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2011 OpenWrt.org
# This is free software, licensed under the GNU General Public License v2.
# See /LICENSE for more information.
#


START=99


start() {
	/usr/bin/oscam -b -u
	echo "runing!"
}

stop() {
	killall oscam
	killall -9 oscam 
	echo "stop!"
}


